# nfs/auto_share.service
[Unit]
Description=NFS automount

[Service]
User=root
Type=simple
ExecStart=/usr/local/bin/auto_share

[Install]
WantedBy=multi-user.target
