# nfs/auto_share.timer
[Unit]
Description=Run auto_share.service every minute

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timer.target

